# apps-expose/com.company.parallels-apps-expose-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.company.parallels-apps-expose-service</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/parallels_apps_expose_service.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/parallels_apps_expose_service.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/parallels_apps_expose_service.stderr</string>
</dict>
</plist>
